# server/sentinelx-server.service
[Unit]
Description=SentinelX Security Monitoring Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=sentinelx
Group=sentinelx
WorkingDirectory=/opt/sentinelx
Environment=CONFIG_PATH=/etc/sentinelx/config.yaml
Environment=LOG_LEVEL=info

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/sentinelx /var/log/sentinelx
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
CapabilityBoundingSet=
LockPersonality=true
MemoryDenyWriteExecute=true

# 执行命令
ExecStart=/opt/sentinelx/bin/sentinelx-server --config /etc/sentinelx/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=10

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinelx

# 限制资源
LimitNOFILE=65536
LimitNPROC=4096
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target